# deploy/systemd/polydata-agent-market-wide-seed.service
[Unit]
Description=polyData market-wide Agent snapshot seed
After=network-online.target redis.service postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/__POLYDATA_REPO_ROOT__
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/bin/bash -lc 'exec "${POLYDATA_PYTHON_BIN:-python3}" scripts/agent/seed_market_wide_insights.py --lens "${POLYDATA_AGENT_MARKET_WIDE_TIMER_LENS:-overview}"'
StandardOutput=journal
StandardError=journal
