# deploy/systemd/polydata-agent-gateway-tunnel.service
[Unit]
Description=polyData reverse SSH tunnel for local AI agent gateway
After=network-online.target polydata-agent-gateway.service
Wants=network-online.target polydata-agent-gateway.service

[Service]
Type=simple
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/bin/bash -lc 'exec ssh -N -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -R "127.0.0.1:${POLYDATA_AGENT_GATEWAY_REMOTE_PORT:-18700}:127.0.0.1:${POLYDATA_AGENT_GATEWAY_PORT:-18700}" "${POLYDATA_AGENT_GATEWAY_REMOTE}"'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

