# deploy/systemd/polydata-polymarket-macro-map-seed.service
[Unit]
Description=polyData Polymarket macro market map seed watcher
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/__POLYDATA_REPO_ROOT__
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/bin/bash -lc 'exec "${POLYDATA_PYTHON_BIN:-python3}" scripts/runtime/polymarket_macro_map_watcher.py --watch --interval "${POLYDATA_MACRO_MARKET_MAP_WATCH_INTERVAL_SECONDS:-180}"'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=polydata-gcp.target
