# deploy/systemd/polydata.target
[Unit]
Description=polyData local production services
Wants=polydata-api.service polydata-market-sync.service polydata-trade-sync.service polydata-oracle-sync.service polydata-analytics-sync.service polydata-new-market-signal.service polydata-geo-sanctions-shock.service polydata-jin10-seed.service polydata-f1-seed.service polydata-nba-seed.service polydata-market-group-seed.service polydata-crypto-funding-seed.service polydata-inflation-nowcast-seed.service polydata-polymarket-macro-map-seed.service polydata-cpi-release-calendar-seed.service polydata-energy-gasoline-shock-seed.service polydata-global-weather-map-seed.service polydata-weather-news-seed.service polydata-food-retail-basket-seed.service polydata-macro-cpi-panels-seed.service polydata-macro-cpi-registry-seed.service polydata-alpha-signal-seed.service polydata-whale-trades-seed.service polydata-suspicious-trades-seed.service polydata-bootstrap-seed.service
After=network-online.target

[Install]
WantedBy=default.target
