# deploy/systemd/polydata-analytics-sync.service
[Unit]
Description=polyData trade analytics sync service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/__POLYDATA_REPO_ROOT__
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/usr/bin/python3 scripts/db/sync_trade_analytics.py --batch-size 200000 --sync-state-key trade_analytics_live --watch --interval 15
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=polydata.target
