# deploy/systemd/polydata-defi-token-watch-seed.service
[Unit]
Description=polyData DeFi token watch seed watcher
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/__POLYDATA_REPO_ROOT__
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/bin/bash -lc 'exec "${POLYDATA_PYTHON_BIN:-python3}" scripts/runtime/defi_token_watch_watcher.py --watch --interval "${POLYDATA_DEFI_TOKEN_WATCH_INTERVAL_SECONDS:-120}" --limit "${POLYDATA_DEFI_TOKEN_WATCH_LIMIT:-10}"'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=polydata-gcp.target
