# deploy/systemd/polydata.target
[Unit]
Description=polyData local production services
Wants=polydata-api.service polydata-market-sync.service polydata-trade-sync.service polydata-oracle-sync.service polydata-analytics-sync.service polydata-new-market-signal.service polydata-geo-sanctions-shock.service
After=network-online.target

[Install]
WantedBy=default.target
