# deploy/systemd/polydata-local-collector.target
[Unit]
Description=polyData local market/orderfilled/oracle collector services
Wants=polydata-market-sync.service polydata-oracle-sync.service polydata-trade-sync.service polydata-analytics-sync.service polydata-event-market-serving.service polydata-db-reverse-tunnel.service
After=network-online.target

[Install]
WantedBy=default.target
