# deploy/systemd/polydata-sports-odds-seed.service
[Unit]
Description=polyData sports odds seed watcher
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/__POLYDATA_REPO_ROOT__
EnvironmentFile=%h/.config/polydata/polydata.env
ExecStart=/bin/bash -lc 'exec "${POLYDATA_PYTHON_BIN:-python3}" scripts/runtime/sports_odds_watcher.py --watch --interval "${POLYDATA_SPORTS_ODDS_WATCH_INTERVAL_SECONDS:-180}" --limit "${POLYDATA_SPORTS_ODDS_LIMIT:-8}"'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=polydata-gcp.target
